# etc/crontab
* 3 * * * /sbin/shutdown -r +5  # root's cron

* 3 * * * /bin/echo " " > /home/chuck/urban/ts/q3ut4/games.log  # chuck's cron
